# platform/systemd/homedir-sdlc-worker.service
[Unit]
Description=HomeDir autonomous SDLC SCC worker (user)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment="HOMEDIR_SDLC_ENV_FILE=/home/homedir-sdlc/.config/homedir-sdlc/env"
ExecStart=/home/homedir-sdlc/.local/bin/homedir-sdlc-worker.sh
Nice=5
NoNewPrivileges=true
PrivateTmp=yes

[Install]
WantedBy=default.target
